# deploy/find-a-home-daily-report.service.example
[Unit]
Description=Find A Home daily tracker-health report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/thibault/find-a-home
EnvironmentFile=/home/thibault/find-a-home/.env
ExecStart=/home/thibault/find-a-home/venv/bin/python -m rental_alert_bot daily-report --config /home/thibault/find-a-home/config.json
