# deploy/find-a-home.service.example
[Unit]
Description=Find A Home rental alert bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/thibault/find-a-home
EnvironmentFile=/home/thibault/find-a-home/.env
ExecStart=/home/thibault/find-a-home/venv/bin/python -m rental_alert_bot run --config /home/thibault/find-a-home/config.json
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
